# root-image/rc.local
#!/bin/bash

echo "Wait for network and devices to settle"
i=0
while [ \$i -lt 60 ]
do
    [ \$(nmcli networking connectivity) == "full" ] && break
    sleep 1
    i=\$(( \$i + 1 ))
done

mount /dev/sr0 /cdrom

PATH=/cdrom:\$PATH
cd /cdrom
./rc || echo "TEST FAILED"
exit 0
